# zfs-remotebackup.cron
2 * * * * root /usr/local/bin/zfs-remotebackup-cron.sh
10 * * * * root FORCE_SNAPTYPE="frequent" /usr/local/bin/zfs-remotebackup-cron.sh
